# rpms/preremove
$!/bin/bash

SERVICE=api
/usr/bin/systemctl stop pepe@${SERVICE}
/usr/bin/systemctl disable pepe@${SERVICE}

#EOF